# mksnaps.target
[Unit]
Description=Make Snapshots
Requires=default.target
